# rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -e

. /usr/bin/common.sh

log_info 'Starting cron ...'
exec s6-setuidgid "$PUID:$PGID" cron